# supervisord.conf
[supervisord]
nodaemon=true

[program:node]
command=/usr/bin/node /var/www/index.js
